# Local/Laravel/PHP 8 & HTTPS & Supervisor & Websockets/.docker/supervisor/conf.d/supervisor-services.conf
[supervisord]
logfile=/etc/supervisor/logs/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=5MB         ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false               ; start in foreground if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[inet_http_server]
port = 9001
username = user # Basic auth username
password = pass # Basic auth password

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:websockets-local]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
user=www-data
stdout_logfile=/var/www/html/storage/logs/laravel.log

[program:queue-worker-local]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --tries=3
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=5
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/queue.log
stopwaitsecs=3600

; [program:queue-worker-staging]
; process_name=%(program_name)s_%(process_num)02d
; command=php /var/www/html/artisan queue:work --tries=3
; autostart=true
; autorestart=true
; stopasgroup=true
; killasgroup=true
; user=www-data
; numprocs=5
; redirect_stderr=true
; stdout_logfile=/var/www/html/storage/logs/queue.log
; stopwaitsecs=3600

; [program:queue-worker-production]
; process_name=%(program_name)s_%(process_num)02d
; command=php /var/www/html/artisan queue:work --tries=3
; autostart=true
; autorestart=true
; stopasgroup=true
; killasgroup=true
; user=www-data
; numprocs=5
; redirect_stderr=true
; stdout_logfile=/var/www/html/storage/logs/queue.log
; stopwaitsecs=3600